# deploy/crontab
# Noda 数据库备份系统 - Crontab 配置
# 格式: 分 时 日 月 周 命令

# 每天凌晨 3:00 执行备份
0 3 * * * /app/backup/backup-postgres.sh >> /var/log/noda-backup/backup.log 2>&1

# 每周日凌晨 3:00 执行验证测试
0 3 * * 0 /app/backup/test-verify-weekly.sh >> /var/log/noda-backup/test.log 2>&1

# 每 6 小时清理旧历史记录
0 */6 * * * /app/backup/lib/metrics.sh cleanup 2>/dev/null || true

# 每天凌晨 4:00 清理旧备份文件
0 4 * * * find /tmp/postgres_backups -type f -name "*.dump" -mtime +7 -delete 2>/dev/null || true

# 每天凌晨 3:30 执行 Doppler 密钥备份（错开数据库备份 3:00）
30 3 * * * /app/backup/backup-doppler-secrets.sh --project noda --config prd >> /var/log/noda-backup/doppler-backup.log 2>&1
