# deploy/supervisord.conf
# ============================================
# Noda Ops - Supervisord 配置
# ============================================
# 管理多个运维进程：备份 cron + Cloudflare Tunnel
# ============================================

[supervisord]
nodaemon=true
logfile=/var/log/noda-backup/supervisord.log
pidfile=/tmp/supervisor/supervisord.pid
childlogdir=/var/log/noda-backup

[program:cron]
command=/usr/sbin/crond -f -l 2
user=root
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/noda-backup/cron.log
stderr_logfile=/var/log/noda-backup/cron-error.log
environment=HOME="/home/nodaops",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",TZ="Pacific/Auckland"

[program:cloudflared]
command=/usr/local/bin/cloudflared tunnel --no-autoupdate run --token %(ENV_CLOUDFLARE_TUNNEL_TOKEN)s
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/var/log/noda-backup/cloudflared.log
stderr_logfile=/var/log/noda-backup/cloudflared-error.log
environment=HOME="/home/nodaops"
# 如果没有配置 CLOUDFLARE_TUNNEL_TOKEN，则不启动
# 通过环境变量控制是否启用
